# container/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:crond]
command=/opt/youtrack/bin/youtrack.sh start